# rootfs/etc/services.d/caddy/run
#!/command/with-contenv bashio
# -*- bash -*-
# shellcheck shell=bash
exec caddy run --config /etc/caddy/Caddyfile