# rootfs/etc/s6-overlay/s6-rc.d/caddy/finish
#!/command/with-contenv bashio
# -*- bash -*-
# shellcheck shell=bash

bashio::log.info "Caddy stopped"