# rootfs/etc/s6-overlay/s6-rc.d/caddy/run
#!/command/with-contenv bashio
# -*- bash -*-
# shellcheck shell=bash

bashio::log.info "Starting Caddy..."
exec caddy run --config /etc/caddy/Caddyfile --adapter caddyfile